# MCserver/assets/openrc_service
#!/sbin/openrc-run

depend() {
    need net
    use logger dns
}

supervisor="screen"
name="minecraft"
command="/usr/bin/screen"
command_args="-DmS minecraft /home/minecraft/server/run-minecraft.sh"
command_user="{$USER}:{$USER}"
pidfile="/var/run/${name}.pid"

start_pre() {
    checkpath --directory --mode 0755 --owner {$USER}:{$USER} /var/run
}

stop_post() {
    su - minecraft -c "screen -S minecraft -p 0 -X stuff 'stop$(printf '\r')'"
}